# deploy/whatsapp-agent-bridge.service
[Unit]
Description=WhatsApp -> Claude Agent bridge
After=network-online.target

[Service]
Type=simple
User=bridge
WorkingDirectory=/opt/whatsapp-agent-bridge
EnvironmentFile=/opt/whatsapp-agent-bridge/.env
ExecStart=/usr/bin/npx tsx src/server.ts
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
